# data/mkp224o.service.template
[Unit]
Description=mkp224o Script - Located at $startshLocation

[Service]
ExecStart=bash $startshLocation

[Install]
WantedBy=multi-user.target